# playbooks/files/cron/archive-notion/crontab
0 0 * * * /home/ansible-user/cron/archive-notion/run.sh
